# services/ngrok.service
[Unit]
Description=ngrok Service
After=network.target

[Service]
EnvironmentFile=/etc/default/webhook_env
ExecStart=/usr/local/bin/ngrok http 5000
WorkingDirectory=${PROJECT_PATH}
Restart=always
User=%u
Group=%u
StandardOutput=append:${PROJECT_PATH}/logs/ngrok.log
StandardError=append:${PROJECT_PATH}/logs/ngrok_error.log

[Install]
WantedBy=multi-user.target

