# services/webhook.service
[Unit]
Description=Webhook Listener Service
After=network.target

[Service]
EnvironmentFile=/etc/default/webhook_env
ExecStart=/usr/bin/python3 ${PROJECT_PATH}/webhook_listener.py
WorkingDirectory=${PROJECT_PATH}
Restart=always
User=%u
Group=%u
StandardOutput=append:${PROJECT_PATH}/logs/webhook.log
StandardError=append:${PROJECT_PATH}/logs/webhook_error.log

[Install]
WantedBy=multi-user.target

